# vagrant/cookbooks/supervisord/templates/default/program.conf.erb
[program:<%= @name %>]
command=<%= @command %>
user=<%= @user %>
directroy=<%= @directroy %>
stdout_logfile=/var/log/<%= @name %>.log
stderr_logfile=/var/log/<%= @name %>.log
