# vagrant/cookbooks/lxc/templates/default/lxc-net.conf.erb
description "lxc network"
author "Serge Hallyn <serge.hallyn@canonical.com>"

start on starting lxc
stop on stopped lxc

env USE_LXC_BRIDGE="false"
env LXC_BRIDGE="lxcbr0"
env LXC_ADDR="<%= node[:lxc][:network][:bridge_ip] %>"
env LXC_NETMASK="255.255.248.0"
env LXC_NETWORK="<%= node[:lxc][:network][:netmask] %>"
env LXC_DHCP_RANGE="<%= node[:lxc][:network][:dhcp_start] %>,<%= node[:lxc][:network][:dhcp_end] %>"
env LXC_DHCP_MAX="2046"
env varrun="/var/run/lxc"

pre-start script
        [ -f /etc/default/lxc ] && . /etc/default/lxc

        [ "x$USE_LXC_BRIDGE" = "xtrue" ] || { stop; exit 0; }

        cleanup() {
                # dnsmasq failed to start, clean up the bridge
                iptables -t nat -D POSTROUTING -s ${LXC_NETWORK} ! -d ${LXC_NETWORK} -j MASQUERADE || true
                ifconfig ${LXC_BRIDGE} down || true
                brctl delbr ${LXC_BRIDGE} || true
        }

        if [ -d /sys/class/net/${LXC_BRIDGE} ]; then
                if [ ! -f ${varrun}/network_up ]; then
                # bridge exists, but we didn't start it
                        stop;
                fi
                exit 0;
        fi

        # set up the lxc network
        echo 1 > /proc/sys/net/ipv4/ip_forward
        mkdir -p ${varrun}
        brctl addbr ${LXC_BRIDGE}
        ifconfig ${LXC_BRIDGE} ${LXC_ADDR} netmask ${LXC_NETMASK} up
        iptables -t nat -A POSTROUTING -s ${LXC_NETWORK} ! -d ${LXC_NETWORK} -j MASQUERADE
        dnsmasq -u lxc-dnsmasq --strict-order --bind-interfaces --pid-file=${varrun}/dnsmasq.pid --conf-file= --listen-address ${LXC_ADDR} --dhcp-range ${LXC_DHCP_RANGE} --dhcp-lease-max=${LXC_DHCP_MAX} --dhcp-no-override --except-interface=lo --interface=${LXC_BRIDGE} || cleanup
        touch ${varrun}/network_up
end script

post-stop script
        [ -f /etc/default/lxc ] && . /etc/default/lxc
        [ -f "${varrun}/network_up" ] || exit 0;
        # if $LXC_BRIDGE has attached interfaces, don't shut it down
        ls /sys/class/net/${LXC_BRIDGE}/brif/* > /dev/null 2>&1 && exit 0;

        if [ -d /sys/class/net/${LXC_BRIDGE} ]; then
                ifconfig ${LXC_BRIDGE} down
                iptables -t nat -D POSTROUTING -s ${LXC_NETWORK} ! -d ${LXC_NETWORK} -j MASQUERADE || true
                pid=`cat ${varrun}/dnsmasq.pid 2>/dev/null` && kill -9 $pid || true
                rm -f ${varrun}/dnsmasq.pid
                brctl delbr ${LXC_BRIDGE}
        fi
        rm -f ${varrun}/network_up
end script