# go/Godeps/_workspace/src/koding/kites/kontrol/files/kontrol.conf
# Ubuntu upstart file at /etc/init/kontrol.conf
# After putting it into /etc/init do the following for initialization:
#
#   initctl reload-configuration
#   initctl list
#
# you should see kontrol in the list, if not you upstart script is wrong. After that start it:
#
#   start kontrol
#
# log is stored into local syslog and /var/log/kontrol.log

description     "Koding Kontrol"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

env KITE_HOME=/etc/kite
env KITE_LOG_NOCOLOR=true
env KONTROL_CONFIG_FILE=/opt/kite/kontrol/config/{{.Environment}}.toml

script
  ulimit -n 100000
  cd /opt/kite/kontrol
  ./kontrol
end script

respawn
respawn limit 3 10

# give up if I respawn 3 times in 60 seconds...

