# vagrant/cookbooks/kd_deploy/templates/default/program.conf.erb
[program:<%= @name %>]
command=<%= @command %>
user=<%= @user %>
directory=<%= @directroy %>
environment=HOME=<%= @directroy %>
stdout_logfile=/var/log/<%= @name %>.log
stderr_logfile=/var/log/<%= @name %>.log
autostart=true
startretries=5
killasgroup=true
