# builders/puppet/modules/hosting_configs/files/etc/profile.d/java.sh
export _JAVA_OPTIONS="-Xmx64m"
export JAVA_HOME="/usr/lib/jvm/jre-1.6.0-sun.x86_64/"

