# server/supervisord.conf
[program:webServer]
command = node index.js
directory = /home/dotcloud/current/
