# go/src/koding/kites/kontrol/files/kontrol.conf
# Ubuntu upstart file at /etc/init/kontrol.conf
# After putting it into /etc/init do the following for initialization:
#
#   initctl reload-configuration
#   initctl list
#
# you should see kontrol in the list, if not you upstart script is wrong. After that start it:
#
#   start kontrol
#
# log is stored into local syslog and /var/log/kontrol.log

description     "koding kontrol"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

script
  cd /opt/kite/kontrol
  ./kontrol -c {{.Profile}} -r {{.Region}} >> /var/log/kontrol.log
end script

respawn
respawn limit 3 10

# give up if I respawn 3 times in 60 seconds...

